# direct/tests-status/etc/cron.d/yandex-du-tests-status
MAILTO = ppc-admin@yandex-team.ru

10 2 * * * root find /var/spool/tests-status -mindepth 2 -maxdepth 2 -mtime +14 -print -delete

