# mail/testing/caldav (2).conf
[program:caldav]
command=java
    -XX:+PrintGCDetails
    -Xloggc:/var/log/calendar/calendar.caldav.gc.log
    -Xms3g
    -Xmx5g
    -showversion
    -server
    -XX:+HeapDumpOnOutOfMemoryError
    -XX:HeapDumpPath=/var/heapdumps
    -Dsun.net.inetaddr.ttl=60
    -Dhttp.maxConnections=100
    -Djava.net.preferIPv6Addresses=true
    -Djava.net.preferIPv4Stack=false
    -Dhttp.keepAlive=true
    -Dfile.encoding=UTF-8
    -Diiop.publishIP=true
    -Detc.dir=/etc/yandex/calendar
    -Dauth.domains=yt
    -Dmicronaut.config.files=/etc/yandex/calendar-micro/application-yt_testing.yml,/etc/yandex/calendar-micro/application.yml
    -Djavax.xml.parsers.DocumentBuilderFactory=org.apache.xerces.jaxp.DocumentBuilderFactoryImpl
    -Djavax.xml.xpath.XPathFactory:http://java.sun.com/jaxp/xpath/dom=com.saxonica.config.ProfessionalXPathFactory
    -Djavax.xml.stream.XMLInputFactory=ru.yandex.calendar.frontend.ews.xml.NullReplacingInputFactory
    -Dorg.joda.time.DateTimeZone.Provider=tzdata.Provider
    -Dlog4j.formatMsgNoLookups=true
    -XX:+UseParallelGC
    -Xdebug
    -XX:ActiveProcessorCount=%(ENV_ACTIVE_PROCESSOR_COUNT)s
    -javaagent:/usr/lib/yandex/calendar/lib/ea-async.jar
    -Xrunjdwp:transport=dt_socket,server=y,suspend=n,address=*:22321
    -classpath /usr/lib/yandex/calendar/lib/calendar.jar:/usr/lib/yandex/calendar/lib/calendar_micro.jar:/usr/lib/yandex/calendar/lib/*:/etc/yandex/logconfig/
    ru.yandex.calendar.boot.Main
    caldav
