# disk/disk-testing/roots/units/logrotate-random/files/etc/cron.d/logrotate-hourly
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

25,55 * * * * root sleep $((RANDOM\%600)) ; [ -x /usr/sbin/logrotate ] && lockf -st0 /tmp/logrotate.lock -- /usr/sbin/logrotate -v /etc/logrotate.conf >> /var/log/logrotate/logrotate.log 2>&1
