# market/mj-test-service/supervisor.conf
[program:mj-test-service]
process_name = %(program_name)s_%(process_num)02d
numprocs = 1
command = /start-service.sh
user = root
autostart = yes
autorestart = yes
startsecs = 2
stopsignal = TERM
stopwaitsecs = 30
directory = %(ENV_PWD)s
environment = LD_PRELOAD="/usr/lib/x86_64-linux-gnu/libjemalloc.so.2",MALLOC_CONF="prof:false,lg_prof_sample:14,prof_prefix:/tmp/jeprof-mj-test-service"
stdout_logfile = /var/log/supervisor/mj-test-service.log
stderr_logfile = /var/log/supervisor/mj-test-service.log
