# mail/deb-package-testing/deploy/etc/supervisor/conf.d/nwsmtp.conf
[program:nwsmtp]
command=/usr/sbin/nwsmtp /etc/nwsmtp/config.yml
priority=15
autostart=true
autorestart=true
startsecs=1
startretries=3
