# mail/deb-package-testing/deploy/etc/supervisor/conf.d/rsyslogd.conf
[program:syslog]
command=/usr/sbin/rsyslogd
priority=10
autostart=true
autorestart=true
