# direct/yandex-dau-secrets-direct-test/etc/cron.d/yandex-dau-secrets-direct-test
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
7-59/10 * * * * root yav-deploy --rsa-login robot-direct-yav-np --rsa-private-key /etc/direct-tokens/ssh-rsa_robot-direct-yav-np
