# disk/disk-testing/roots/units/logbackup/files/etc/cron.d/disk-logbackup
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
#MAILTO=disk-root@yandex-team.ru

30 4 * * * root flock -w 600 /var/run/disk-logbackup.lock /usr/lib/yandex/disk/logbackup/disk-logbackup.sh >> /var/log/disk-logbackup.log 2>&1
