# disk/disk-testing/roots/units/sophos/files/etc/cron.d/sophos
10  */6 *   *   *   root    /usr/sbin/sophos-update.sh > /var/log/sophos-update.log 2>&1
*/2 *   *   *   *   root    /usr/bin/savdid-watchdog.sh
