# disk/disk-testing/roots/units/downloader-counter/files/etc/cron.d/yandex-disk-downloader-counter
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
#MAILTO=disk-root@yandex-team.ru

* * * * * root flock -x -w 5 /var/lock/disk_downloader.count-reporter.lock /usr/bin/disk_downloader.count-reporter.pl > /dev/null
