# disk/disk-testing/roots/units/esets/files/etc/cron.d/esets
*   *   *   *   *   root    /usr/bin/esets-watchdog.sh
