# disk/disk-testing/roots/units/coredumps-cleaner/files/etc/cron.d/coredumps-cleaner
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
#MAILTO=disk-root@yandex-team.ru

05 4 * * *  root    flock -n /tmp/coredumps-cleaner.lock /usr/bin/coredumps-cleaner.sh 30 &>/dev/null
