# mail/deb-package-testing/deploy/etc/supervisor/conf-available/notsolitesrv.conf
[program:notsolitesrv]
command=/usr/bin/notsolitesrv /etc/notsolitesrv/config.yml
priority=10
autostart=true
autorestart=true
