# mail/deb-package-testing/deploy/etc/supervisor/conf-available/nwsmtp-out.conf
[program:nwsmtp-out]
command=/usr/sbin/nwsmtp /etc/nwsmtp/config-out.yml
priority=15
autostart=true
autorestart=true
startsecs=1
startretries=3

