# mail/deb-package-testing/deploy/etc/supervisor/conf-available/smtpgate.conf
[program:smtpgate]
command=/usr/sbin/smtpgate /etc/smtpgate/smtpgate.conf.yml
priority=10
autostart=true
autorestart=true
