# market/GENERAL/yandex-proftpd4-fw-testing.conf
description "Market2Sber proxy"

start on (net-device-up
          and local-filesystems
          and static-network-up
          and runlevel [2345])

stop on runlevel [016]

exec /usr/sbin/proftpd -c /etc/yandex/proftpd/proftpd4-fw-testing.conf -n -d 10
respawn
