# disk/disk-testing/roots/units/yarl/files/etc/cron.d/yarl-update-root
SHELL=/bin/sh
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin

*/1 * * * * root flock -w 30 /tmp/yarl-update-root.lock /usr/local/bin/yarl-update-root.sh
