# travel/testcontext/supervisor/10-app.conf
[program:testcontext]
directory = /app/testcontext/
command = /app/testcontext/testcontext
environment=CONFIG_PATH="/app/testcontext/config.%(ENV_YENV_TYPE)s.yaml"

autostart = true
autorestart = true
startsecs = 5
startretries = 1000
stopwaitsecs = 30

stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes = 0
