# disk/disk-testing/roots/units/yandex-disk-downloader/files/etc/cron.d/yandex-disk-downloader-meta
PATH=/bin:/sbin:/usr/bin:/usr/sbin:/usr/local/bin:/usr/local/sbin
#MAILTO=disk-root@yandex-team.ru

*/15 * * * * root flock -n /tmp/rm-nginx-spool.lock /usr/bin/disk_downloader.clean-tmp.sh

