# disk/disk-testing/roots/units/yandex-selfdns-client/files/etc/cron.d/yandex-selfdns-client
#For testing
#*/5 * * * * selfdns /usr/bin/selfdns-client
#For production
*/5 * * * *     selfdns     sleep $((RANDOM\%60+60)) && /usr/bin/selfdns-client 1>>/var/log/yandex-selfdns-client/std.log 2>&1

